# Palo-Alto-EDL/pan-edl-www.service
[Unit]
Description=Simple Python HTTP server serving text file for Pan EDL
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/pan-edl-list
ExecStart=/usr/bin/python3 -m http.server 8080 --bind 0.0.0.0
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
